# scripts/lcloud-backup.service
[Unit]
Description=LCloud daily backup
Documentation=file:///root/LCloud/scripts/backup.sh
After=network.target

[Service]
Type=oneshot
ExecStart=/root/LCloud/scripts/backup.sh
StandardOutput=append:/root/LCloud/logs/backup.log
StandardError=append:/root/LCloud/logs/backup.log
Nice=10
IOSchedulingClass=idle
